# deploy.systemd.service.example
[Unit]
Description=Python Web UI Request Tester
After=network.target

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/opt/python-webui-3-token-import-argparse
Environment=PYTHONUNBUFFERED=1
Environment=REQUEST_TESTER_HOST=127.0.0.1
Environment=REQUEST_TESTER_PORT=8765
Environment=REQUEST_TESTER_DATA_DIR=/var/lib/python-webui-request-tester
StateDirectory=python-webui-request-tester
ExecStart=/opt/python-webui-3-token-import-argparse/.venv/bin/python server.py
Restart=on-failure
RestartSec=3

[Install]
WantedBy=multi-user.target
